# pac-n90.service
[Unit]
Description=PAC N90 MQTT service.

[Service]
Type=simple
ExecStart=/usr/local/bin/pac-n90

[Install]
WantedBy=multi-user.target
